# services/local-rerank-sidecar/com.openclaw.local-rerank-sidecar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.openclaw.local-rerank-sidecar</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/lucifinil_chen/.local/bin/uv</string>
      <string>run</string>
      <string>serve.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/lucifinil_chen/.openclaw/workspace/services/local-rerank-sidecar</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>LOCAL_RERANK_HOST</key>
      <string>127.0.0.1</string>
      <key>LOCAL_RERANK_PORT</key>
      <string>8765</string>
      <key>LOCAL_RERANK_MODEL</key>
      <string>BAAI/bge-reranker-v2-m3</string>
      <key>LOCAL_RERANK_LOG_LEVEL</key>
      <string>INFO</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/lucifinil_chen/.openclaw/workspace/runtime/local-rerank-sidecar/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lucifinil_chen/.openclaw/workspace/runtime/local-rerank-sidecar/launchd.log</string>
  </dict>
</plist>
